# src/install/tilt-machine.service
[Unit]
Description=Start the Monitor Rotation Service

[Service]
ExecStart=
Type=oneshot
